# immurok-daemon.service
[Unit]
Description=immurok BLE fingerprint authentication daemon
After=bluetooth.target dbus.service

[Service]
Type=simple
ExecStart=%h/.local/bin/immurok-daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
